# files/etc/init/edgecore.conf
description "KubeEdge - edgecore"

start on started containerd
stop on runlevel [!2345]

respawn

kill timeout 20

script
	EDGECORE=/usr/local/bin/edgecore
	exec "$EDGECORE"
end script
